# hellofamily_celery.conf
[program:celery]
command=celery -A app.tasks worker --loglevel info --logfile celery_worker.log
directory=/var/www/hellofamily

numprocs=1
autostart=true
autorestart=true